# services/teaspeak_helper/run
#!/bin/bash
# TeaSpeak helper service - runs alongside the main process
# For now, just sleep indefinitely to keep the service running
while true; do
    sleep 3600
done