# services/teaspeak_helper/run
#!/bin/bash
exec s6-setuidgid tea /opt/teaspeak/scripts/helper.sh