# services/teaspeak/finish
#!/command/with-contenv sh
echo "TeaSpeak service finished"