# services/teaspeak/run
#!/bin/bash

# Source environment variables
source /etc/environment

# Change to TeaSpeak directory
cd /opt/teaspeak

# Execute the startup.sh script
if [ -f "/opt/teaspeak/scripts/startup.sh" ]; then
    exec /opt/teaspeak/scripts/startup.sh
else
    echo "Error: startup.sh not found!"
    exit 1
fi